# deploy/vakit-pi.service
# Vakit-Pi Systemd Service Template
# ===================================
# NOT: Bu dosya bir şablondur!
# install.sh scripti bu dosyayı otomatik olarak doğru yollarla oluşturur.
# Manuel kurulum yapıyorsanız, aşağıdaki değerleri kendi sisteminize göre düzenleyin:
#   - User/Group: Kullanıcı adınız (örn: pi)
#   - WorkingDirectory: Proje dizini (örn: /home/efraim/vakit-pi)
#   - VAKIT_PI_SETTINGS_PATH: Ayar dosyası yolu
#   - ExecStart: uv yolu ve proje dizini
#   - ReadWritePaths: Ayar dizini

[Unit]
Description=Vakit-Pi - Namaz Vakti ve Ezan Uygulaması
Documentation=https://github.com/vakit-pi/vakit-pi
After=network.target sound.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=efraim
Group=efraim
WorkingDirectory=/home/efraim/vakit-pi

# Environment variables
Environment="VAKIT_PI_HOST=0.0.0.0"
Environment="VAKIT_PI_PORT=8080"
Environment="VAKIT_PI_LOG_LEVEL=INFO"
Environment="VAKIT_PI_SETTINGS_PATH=/home/efraim/.config/vakit-pi/settings.json"

# Uygulama komutu (uv kullanarak)
ExecStart=/home/efraim/.local/bin/uv run vakit-pi serve

# Alternatif: pip ile kurulduysa
# ExecStart=/home/efraim/.local/bin/vakit-pi serve

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
# ProtectHome devre dışı - uv cache erişimi için gerekli

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vakit-pi

[Install]
WantedBy=multi-user.target
